# systemd/dotfiles-autocommit.timer
[Unit]
Description=Daily dotfiles auto-commit

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=10m

[Install]
WantedBy=timers.target
